# app/rootfs/etc/s6/services/mail/run
#!/bin/sh

exec tail -f --retry /var/mail/mail
