# docker/ldap/rootfs/etc/s6/services/slapd/run
#!/bin/sh

slapd-fix
exec slapd -d 0 -h "ldap:/// ldapi:///" -u openldap -g openldap ${LDAP_OPTS:-}
