# app/rootfs/etc/s6/services/webldap/run
#!/bin/sh

cd /srv/webldap

webldap-wrap python3 manage.py syncdb
exec webldap-wrap python3 manage.py runserver 0.0.0.0:8000
